# com.robert.dotfiles-pull.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.robert.dotfiles-pull</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/git</string>
        <string>-C</string>
        <string>/Users/Robert/repos/dotfiles</string>
        <string>pull</string>
        <string>--ff-only</string>
        <string>origin</string>
        <string>main</string>
    </array>

    <key>StartInterval</key>
    <integer>604800</integer>

    <key>StandardOutPath</key>
    <string>/tmp/dotfiles-pull.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/dotfiles-pull.log</string>

    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
